# modules/gunicorn/templates/mezzanine_supervisor.conf.erb
[program:myproject]
command = gunicorn --pythonpath /home/sanoma/myproject/ myproject.wsgi
directory = /home/sanoma
user = root 
