# modules/gunicorn/templates/mezzanine_supervisor.conf.erb
[program:<%= @project %>]
command = gunicorn --pythonpath /home/<%= @user %>/<%= @project %>/ <%= @project %>.wsgi
directory = /home/<%= @user %>
user = root 
